# crontab
* * * * *  python /home/app/test_github.py >>/tmp/out.log 2>/tmp/err.log